# crontab
59 23 * * * /root/backup.sh 2>>/root/backup.log
